# MergerFS-Rclone/Upload Scripts/rclone-uploader.timer
[Unit]
Description=Run Rclone Upload daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target